# docker/cron/crontab
# PageSpace Cron Jobs
# All times are in UTC
# Jobs run inside Docker network, hitting internal web service
# Authentication: HMAC-signed requests via cron-curl helper (set up by entrypoint.sh)

# Memory processing - Daily at 6am UTC
# Discovers patterns from user conversations and updates personalization
0 6 * * * cron-curl POST __WEB_URL__/api/memory/cron >> /var/log/cron/memory.log 2>&1

# Pulse summaries - Every 6 hours
# Generates workspace activity summaries for active users
0 */6 * * * cron-curl POST __WEB_URL__/api/pulse/cron >> /var/log/cron/pulse.log 2>&1

# Token cleanup - Hourly
# Removes expired device tokens from database
0 * * * * cron-curl GET __WEB_URL__/api/cron/cleanup-tokens >> /var/log/cron/cleanup.log 2>&1

# Calendar sync - Every 5 minutes
# Syncs Google Calendar connections that are due based on their sync frequency
*/5 * * * * cron-curl GET __WEB_URL__/api/cron/calendar-sync >> /var/log/cron/calendar-sync.log 2>&1

# Retention cleanup - Daily at 1am UTC
# Deletes expired rows from sessions, tokens, page versions, drive backups, permissions, monitoring data
0 1 * * * cron-curl GET __WEB_URL__/api/cron/retention-cleanup >> /var/log/cron/retention-cleanup.log 2>&1

# Security audit chain verification - Daily at 2am UTC
# Verifies hash chain integrity of security audit log
0 2 * * * cron-curl GET __WEB_URL__/api/cron/verify-audit-chain >> /var/log/cron/audit-chain.log 2>&1

# Activity log archival (hot→cold tier) - Daily at 2:30am UTC
# Flips isArchived=true on activity_logs older than ACTIVITY_LOGS_ARCHIVE_DAYS (default 365)
# so the active-view stays bounded. Flag-flip only: never deletes rows or touches the hash chain.
30 2 * * * cron-curl GET __WEB_URL__/api/cron/archive-activity-logs >> /var/log/cron/archive-activity-logs.log 2>&1

# AI usage log purge - Daily at 3am UTC
# Anonymizes prompt/completion text >30 days, purges rows >90 days
0 3 * * * cron-curl POST __WEB_URL__/api/cron/purge-ai-usage-logs >> /var/log/cron/purge-ai-usage.log 2>&1

# Deleted message purge - Daily at 4am UTC
# Hard-deletes soft-deleted messages and conversations older than 30 days
0 4 * * * cron-curl POST __WEB_URL__/api/cron/purge-deleted-messages >> /var/log/cron/purge-messages.log 2>&1

# Workflow execution - Every 5 minutes
# Checks for due scheduled workflows and executes them
*/5 * * * * cron-curl POST __WEB_URL__/api/cron/workflows >> /var/log/cron/workflows.log 2>&1

# Calendar trigger execution - Every 2 minutes
# Fires scheduled agent work tied to calendar events
*/2 * * * * cron-curl POST __WEB_URL__/api/cron/calendar-triggers >> /var/log/cron/calendar-triggers.log 2>&1

# Trashed page purge - Daily at 5am UTC
# Hard-deletes pages that have been in the trash for 30+ days (GDPR Art. 17 #909)
0 5 * * * cron-curl GET __WEB_URL__/api/cron/purge-trashed-pages >> /var/log/cron/purge-trashed-pages.log 2>&1

# Orphaned file cleanup - Weekly on Sundays at 6am UTC
# Detects and deletes file records with zero references + physical files
0 6 * * 0 cron-curl GET __WEB_URL__/api/cron/cleanup-orphaned-files >> /var/log/cron/orphan-cleanup.log 2>&1

# Subscription tier reconcile - Every 30 minutes
# Compares users.subscriptionTier (denormalized cache) against the
# subscriptions table and repairs repairable drift; indeterminate drift
# (an entitled subscription on an unmapped price id) is logged/audited for
# a human, never auto-repaired. Backstop for missed/failed Stripe webhooks
# (see subscription-tier-reconcile.ts in @pagespace/lib).
*/30 * * * * cron-curl GET __WEB_URL__/api/cron/reconcile-subscription-tiers >> /var/log/cron/reconcile-subscription-tiers.log 2>&1

# Credit reconcile - Every 10 minutes
# Re-settles ledger rows stuck 'pending' and bills usage rows that never
# decremented the prepaid balance, so every AI call is charged exactly once
# across crashes/deploys. Local-only — makes no Stripe calls.
*/10 * * * * cron-curl GET __WEB_URL__/api/cron/reconcile-credits >> /var/log/cron/reconcile-credits.log 2>&1

# Machine storage reconcile - Hourly
# Bills each Terminal Machine's persistent storage from its last persisted
# measurement (see machine-storage-reconcile.ts in @pagespace/lib); watermarks
# make sequential reruns and missed runs self-correcting, but CONCURRENT runs
# can double-bill. flock -n is a best-effort guard so this container's own
# ticks never overlap (a skipped tick is caught up by the watermark next
# hour); in-service serialization covering ALL callers is tracked follow-up.
0 * * * * flock -n /run/reconcile-machine-storage.lock cron-curl GET __WEB_URL__/api/cron/reconcile-machine-storage >> /var/log/cron/reconcile-machine-storage.log 2>&1 || echo "$(date -u) tick skipped (lock held) or request failed" >> /var/log/cron/reconcile-machine-storage.log

# Orphaned sprite reconcile - Every 30 minutes
# Destroys Sprites (billing microVMs) that nothing in the product can reach any
# more. Two sources — see machine-orphan-reconcile.ts in @pagespace/lib:
#   (A) machine_sprite_reclaims, the outbox. AFTER DELETE triggers (migration
#       0210) rescue a sandboxId out of machine_sessions/machine_branches as its
#       row is cascaded away by ANY page/drive/user hard delete — the 30-day
#       purge, "delete permanently", account erasure. The pointer outlives the
#       resource, so no delete path can strand a VM, and none of them need a
#       guard (erasure must never be blocked by a Sprite we failed to kill).
#   (B) tracking rows under a trashed page whose teardown was REQUESTED but never
#       confirmed — a deleteMachine whose kill failed.
# A Machine merely dragged to the trash is left ALONE: its Sprite hibernates and
# a restore is expected to hand the disk back (a kill is irreversible, a trash is
# not). Every effect is safe under concurrency (the kill treats a 404 as success;
# the release writes are compare-and-swaps), so unlike the storage reconcile this
# needs no advisory lock — flock -n just keeps this container's own ticks from
# doing redundant work.
*/30 * * * * flock -n /run/reconcile-orphaned-sprites.lock cron-curl GET __WEB_URL__/api/cron/reconcile-orphaned-sprites >> /var/log/cron/reconcile-orphaned-sprites.log 2>&1 || echo "$(date -u) tick skipped (lock held) or request failed" >> /var/log/cron/reconcile-orphaned-sprites.log

# Expired row sweep - Every 15 minutes
# Reaps expired tombstones from append-with-TTL tables (revoked_service_tokens;
# rate_limit_buckets in PR3; ai_pending_abort_intents). Short cadence matches the 5-minute default
# service-token lifetime so the table never grows more than a few cycles deep.
*/15 * * * * cron-curl GET __WEB_URL__/api/cron/sweep-expired >> /var/log/cron/sweep-expired.log 2>&1

# Scheduled drive backups - Hourly
# Fires automated drive backups for pro+ users whose nextRunAt is due
0 * * * * cron-curl POST __WEB_URL__/api/cron/scheduled-backups >> /var/log/cron/scheduled-backups.log 2>&1
