# docker/cron/crontab
# PageSpace Cron Jobs
# All times are in UTC
# Jobs run inside Docker network, hitting internal web service
# Authentication: HMAC-signed requests via cron-curl helper (set up by entrypoint.sh)

# Memory processing - Daily at 6am UTC
# Discovers patterns from user conversations and updates personalization
0 6 * * * cron-curl POST http://web:3000/api/memory/cron >> /var/log/cron/memory.log 2>&1

# Pulse summaries - Every 6 hours
# Generates workspace activity summaries for active users
0 */6 * * * cron-curl POST http://web:3000/api/pulse/cron >> /var/log/cron/pulse.log 2>&1

# Token cleanup - Hourly
# Removes expired device tokens from database
0 * * * * cron-curl GET http://web:3000/api/cron/cleanup-tokens >> /var/log/cron/cleanup.log 2>&1

# Calendar sync - Every 5 minutes
# Syncs Google Calendar connections that are due based on their sync frequency
*/5 * * * * cron-curl GET http://web:3000/api/cron/calendar-sync >> /var/log/cron/calendar-sync.log 2>&1
